# nerdy-print-shop.service
[Unit]
Description=NerdyPrint - Marketplace de Impressão 3D
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nerdy-print-shop
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="VITE_APP_ID=nerdy-print-shop"
Environment="VITE_OAUTH_PORTAL_URL=https://oauth.manus.computer"
Environment="DATABASE_URL=file:/home/ubuntu/nerdy-print-shop/nerdy_print_shop.db"
Environment="JWT_SECRET=your-secure-random-string-here"
Environment="OAUTH_SERVER_URL=https://oauth.manus.computer"
Environment="OWNER_OPEN_ID=mateinorolamento89@gmail.com"
Environment="BUILT_IN_FORGE_API_URL=https://forge.manus.im"
Environment="BUILT_IN_FORGE_API_KEY=sk-xgf_Doi_U_DXMVLWQXW4z7zoHMjfzQpsJlcAQT6GihBf9Cwxv65flznwcEeuA8y1ZSnY6wu-ORCNSg8vVFPIG0rt71UK"
Environment="WHATSAPP_PHONE_NUMBER=+5511953739362"

ExecStart=/usr/bin/node /home/ubuntu/nerdy-print-shop/dist/index.js

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nerdy-print-shop

[Install]
WantedBy=multi-user.target
